# systemd/ai-reminder.service.example
[Unit]
Description=AI Reminder Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ai-reminder
ExecStart=/home/pi/ai-reminder/.venv/bin/python3 /home/pi/ai-reminder/main.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target